# deploy/sotd.service
# Ship of Tools backend daemon — systemd user unit TEMPLATE (ADR 0030 §3).
#
# Installers substitute the @TOKENS@ and place the result at
# ~/.config/systemd/user/sotd.service, then:
#   systemctl --user daemon-reload && systemctl --user enable --now sotd
#   loginctl enable-linger $USER    # keep it alive without an open session
#
# Token values:
#   @SOT_BIN@          dev fleet: <repo>/rust/target/release/sotd
#                      public install: ~/.local/share/sot/bin/sotd (ADR 0030 §4)
#   @SOT_PORT@         backend TCP port; the fleet convention is 18743
#   @SOT_PROJECT_ROOT@ the project the daemon serves
#
# Shape mirrors the proven live unit on the dev backend (enabled + active since the
# ADR-0028 cutover). Loopback bind only — remote FEs come in over SSH
# forwards, never an exposed bind (ADR 0010).

[Unit]
Description=Ship of Tools backend daemon (sotd)
After=default.target network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=@SOT_BIN@ --tcp 127.0.0.1:@SOT_PORT@ --project-root @SOT_PROJECT_ROOT@ --label sot
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
